# playbooks/roles/vhost/templates/etc/systemd/system/sync-logs-on-exit.service.j2
[Unit]
Description=Synchronizes logs on exit
After=rsyslog.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/edx/bin/sync-logs-on-exit

[Install]
WantedBy=multi-user.target
