# playbooks/edx-east/roles/snort/templates/etc/cron.daily/oinkmaster.j2
#! /bin/bash

oinkmaster -C /etc/oinkmaster.conf -o /etc/snort/rules/ > /dev/null
service snort restart
