# playbooks/roles/browsers/templates/xvfb.service.j2
[Unit]
Description=Xvfb X Server
After=network.target

[Service]
ExecStart=/usr/bin/Xvfb {{ browser_xvfb_display }} -screen 0 1024x768x24

[Install]
WantedBy=multi-user.target
