# playbooks/roles/aws/templates/etc/init/sync-on-stop.conf.j2
start on stopped supervisor
description     "sync s3 logs on supervisor shutdown"
script
        /bin/bash {{ aws_s3_sync_script }} 
end script
