# playbooks/roles/common/templates/etc/cron.hourly/logrotate.j2
#!/bin/sh

test -x /usr/sbin/logrotate || exit 0
/usr/sbin/logrotate /etc/logrotate.d/hourly
