# playbooks/roles/edxapp/templates/edx/app/supervisor/conf.d.available/sshd.conf.j2
[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile={{ supervisor_log_dir }}/%(program_name)s-stdout.log
stderr_logfile={{ supervisor_log_dir }}/%(program_name)s-stderr.log
autorestart=true
