# playbooks/roles/jenkins_worker/templates/python_pkgs.conf.j2
# Pre-install Python packages from S3

start on runlevel [2345]

task
exec install_python_pkgs.sh {{ PYTHON_PKG_URL }}
