# playbooks/roles/security/templates/etc/systemd/system/apt-daily-upgrade.timer.d/override.conf
[Timer]
OnCalendar=
Persistent=false
OnBootSec=1d
OnUnitActiveSec=1d
RandomizedDelaySec=60m
