# playbooks/roles/newrelic/templates/etc/init/newrelic-logwatch-agent.conf.j2
description  "newrelic nginx agent"
start on runlevel [2345]
stop on runlevel [016]
respawn
chdir {{ newrelic_logwatch_dir }} 
exec ruby newrelic_logwatcher_agent.rb
