# playbooks/roles/newrelic/templates/etc/init/newrelic-nginx-agent.conf.j2
description  "newrelic nginx agent"
start on runlevel [2345]
stop on runlevel [016]
respawn
exec /opt/newrelic_nginx_agent/newrelic_nginx_agent.daemon run -t
