# playbooks/roles/edxapp/templates/celerybeat.conf.j2
# generated by ansible

{% if EDXAPP_CELERY_BEAT_LMS_CONFIG.enabled %}
[program:lms_celerybeat]

environment={% if COMMON_ENABLE_NEWRELIC_APP %}NEW_RELIC_APP_NAME={{ EDXAPP_NEWRELIC_WORKERS_APPNAME }},NEW_RELIC_LICENSE_KEY={{ NEWRELIC_LICENSE_KEY }},{% endif -%}LOGLEVEL=info,DJANGO_SETTINGS_MODULE={{ worker_django_settings_module }},PYTHONPATH={{ edxapp_code_dir }},SERVICE_VARIANT=lms,BOTO_CONFIG="{{ edxapp_app_dir }}/.boto"
user={{ common_web_user }}
directory={{ edxapp_code_dir }}
stdout_logfile={{ supervisor_log_dir }}/%(program_name)s-stdout.log
stderr_logfile={{ supervisor_log_dir }}/%(program_name)s-stderr.log

command={{ edxapp_venv_dir + '/bin/newrelic-admin run-program ' if EDXAPP_CELERY_BEAT_LMS_CONFIG.monitor and COMMON_ENABLE_NEWRELIC_APP else ''}}{{ edxapp_venv_bin }}/python {{ edxapp_code_dir }}/manage.py lms --settings={{ worker_django_settings_module }} celerybeat --loglevel=info --pidfile={{ EDXAPP_CELERY_BEAT_LMS_CONFIG.process_id }} --schedule {{ EDXAPP_CELERY_BEAT_LMS_CONFIG.schedule_file }}

stopwaitsecs={{ EDXAPP_CELERY_BEAT_LMS_CONFIG.stopwaitsecs | default(EDXAPP_WORKER_DEFAULT_STOPWAITSECS) }}
{% endif -%}
{% if EDXAPP_CELERY_BEAT_CMS_CONFIG.enabled %}
[program:cms_celerybeat]

environment={% if COMMON_ENABLE_NEWRELIC_APP %}NEW_RELIC_APP_NAME={{ EDXAPP_NEWRELIC_WORKERS_APPNAME }},NEW_RELIC_LICENSE_KEY={{ NEWRELIC_LICENSE_KEY }},{% endif -%}LOGLEVEL=info,DJANGO_SETTINGS_MODULE={{ worker_django_settings_module }},PYTHONPATH={{ edxapp_code_dir }},SERVICE_VARIANT=cms,BOTO_CONFIG="{{ edxapp_app_dir }}/.boto"
user={{ common_web_user }}
directory={{ edxapp_code_dir }}
stdout_logfile={{ supervisor_log_dir }}/%(program_name)s-stdout.log
stderr_logfile={{ supervisor_log_dir }}/%(program_name)s-stderr.log

command={{ edxapp_venv_dir + '/bin/newrelic-admin run-program ' if EDXAPP_CELERY_BEAT_LMS_CONFIG.monitor and COMMON_ENABLE_NEWRELIC_APP else ''}}{{ edxapp_venv_bin }}/python {{ edxapp_code_dir }}/manage.py cms --settings={{ worker_django_settings_module }} celerybeat --loglevel=info --pidfile={{ EDXAPP_CELERY_BEAT_LMS_CONFIG.process_id }} --schedule {{ EDXAPP_CELERY_BEAT_LMS_CONFIG.schedule_file }}

stopwaitsecs={{ EDXAPP_CELERY_BEAT_LMS_CONFIG.stopwaitsecs | default(EDXAPP_WORKER_DEFAULT_STOPWAITSECS) }}
{% endif -%}
#[group:edxapp_celerybeat]
#programs=lms,cms
