# playbooks/roles/hotg/templates/edx/app/supervisor/conf.d/hotg.conf.j2
[program:{{ HOTG_SERVICE_NAME  }}]
directory={{ hotg_app_dir }}
command={{ hotg_app_dir }}/supervisor_wrapper.sh
stdout_logfile=syslog
stderr_logfile=syslog
user={{ common_web_user }}